# test.service
[Unit]
Description=monitor test process

[Service]
Type=oneshot
ExecStart=/usr/local/bin/test.sh

